# supervisord.conf
[supervisord]
nodaemon=true                  ; run in foreground (required for Docker)
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ─── MCP Server 1 on port 8000 ───
[program:db-mcp]
command=uvicorn server:app --host 0.0.0.0 --port 8000 --workers 1
directory=/cm-mcp/db-mcp/src
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/db-mcp.err.log
stdout_logfile=/var/log/db-mcp.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"

; ─── MCP Server 2 on port 8001 ───
[program:ask-uos-mcp]
command=uvicorn server:app --host 0.0.0.0 --port 8001 --workers 1
directory=/cm-mcp/ask-uos-mcp
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/ask-uos-mcp.err.log
stdout_logfile=/var/log/ask-uos-mcp.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"

